# roles/backups/files/backups.service
[Unit]
Description=NAS Backups

[Service]
Type=oneshot
ExecStart=/bin/sh -c /usr/local/bin/backups.sh
